# ops/systemd/iic-triage.service
[Unit]
Description=IIC triage consumer — dedupe + salience + persistence
After=network-online.target redis-server.service
Requires=redis-server.service

[Service]
Type=simple
User=iic
WorkingDirectory=/home/iic/TradingAgents/TradingAgents
EnvironmentFile=/home/iic/TradingAgents/TradingAgents/.env
ExecStart=/home/iic/TradingAgents/TradingAgents/.venv/bin/python -m tradingagents.sensing.triage
Restart=on-failure
RestartSec=30
MemoryMax=2G
CPUQuota=100%
StandardOutput=append:/var/log/iic/triage.log
StandardError=append:/var/log/iic/triage.log

[Install]
WantedBy=multi-user.target
