# ops/systemd/iic-promoter.service
# /etc/systemd/system/iic-promoter.service
#
# IIC-FORGE F4 — promoter service
# Polls events table every 10s, enqueues event_alert jobs into queue_jobs.
# Defensive retry-internal: see tradingagents/orchestrator/promoter.py
# The exit-gate evaluator (scripts/f4_exit_gate.py) fails the gate if this
# unit's NRestarts > 0 during the observation window.

[Unit]
Description=IIC-FORGE — F4 promoter (events → queue_jobs)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=iic
WorkingDirectory=/home/iic/TradingAgents/TradingAgents
EnvironmentFile=/home/iic/TradingAgents/TradingAgents/.env
ExecStart=/home/iic/TradingAgents/TradingAgents/.venv/bin/python -m tradingagents.orchestrator.promoter
Restart=on-failure
RestartSec=30
MemoryMax=256M
CPUQuota=25%
StandardOutput=append:/var/log/iic/promoter.log
StandardError=append:/var/log/iic/promoter.log

[Install]
WantedBy=multi-user.target
