# ops/systemd/iic-sense-macro.service
[Unit]
Description=IIC sensing — Macro (FRED + TE) adapter
After=network-online.target redis-server.service
Requires=redis-server.service

[Service]
Type=simple
User=iic
WorkingDirectory=/home/iic/TradingAgents/TradingAgents
EnvironmentFile=/home/iic/TradingAgents/TradingAgents/.env
ExecStart=/home/iic/TradingAgents/TradingAgents/.venv/bin/python -m tradingagents.sensing.adapters.macro
Restart=on-failure
RestartSec=30
MemoryMax=512M
CPUQuota=50%
StandardOutput=append:/var/log/iic/sense-macro.log
StandardError=append:/var/log/iic/sense-macro.log

[Install]
WantedBy=multi-user.target
