# ops/systemd/iic-watchlist-sweep.service
[Unit]
Description=IIC watchlist TTL sweep
After=network-online.target

[Service]
Type=oneshot
User=iic
WorkingDirectory=/home/iic/TradingAgents/TradingAgents
EnvironmentFile=/home/iic/TradingAgents/TradingAgents/.env
ExecStart=/home/iic/TradingAgents/TradingAgents/.venv/bin/python -m cli.main forge sense sweep-watchlist
StandardOutput=append:/var/log/iic/watchlist-sweep.log
StandardError=append:/var/log/iic/watchlist-sweep.log
